# channels.ini
[program:cmdb_ASGI/Channels]
command=/opt/cmdb/bin/python manage.py runserver localhost:8000
directory = /opt/mycmdb
autostart=true
autorestart=true
user=root
stopsignal=QUIT
startsecs = 5
startretries = 3
redirect_stderr=true
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stdout_logfile=/opt/mycmdb/log/channels_stdout.log
stderr_logfile=/opt/mycmdb/log/channels_err.log
