# cmdb.ini
[program:cmdb_uwsgi]
directory = /opt/mycmdb
command=/usr/bin/uwsgi --ini /opt/mycmdb/uwsgi.ini
autostart=true
autorestart=true
user=root
stopsignal=QUIT
startsecs = 5
startretries = 3
redirect_stderr=true
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stdout_logfile=/opt/mycmdb/log/cmdb_stdout.log
stderr_logfile=/opt/mycmdb/log/cmdb_err.log
